# auto_eject_sda.service
[Unit]
Description=Auto Eject sda
After=udev.service

[Service]
Type=oneshot
ExecStart=ejectsda

[Install]
WantedBy=multi-user.target
